# root/etc/s6-overlay/s6-rc.d/init-mod-swag-maxmind-setup/up
/etc/s6-overlay/s6-rc.d/init-mod-swag-maxmind-setup/run