# root/etc/s6-overlay/s6-rc.d/svc-mod-swag-geoip2influx/run
#!/usr/bin/with-contenv bash

exec \
    python3 /geoip2influx/run.py