# root/etc/s6-overlay/s6-rc.d/init-mod-transmission-transmission-web-control-install/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

printf "/themes/transmission-web-control" > /var/run/s6/container_environment/TRANSMISSION_WEB_HOME

lsiown -R abc:abc /themes
