# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-svn-add-package/run
#!/usr/bin/with-contenv bash

if ! dpkg -s subversion >/dev/null 2>&1; then
    echo "**** Adding subversion to package install list ****"
    echo "subversion" >> /mod-repo-packages-to-install.list
else
    echo "**** subversion already installed, skipping ****"
fi
