# root/etc/s6-overlay/s6-rc.d/init-mod-lidarr-flac2mp3-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-lidarr-flac2mp3-add-package/run