# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-nvm-install/up
/etc/s6-overlay/s6-rc.d/init-mod-code-server-nvm-install/run