# root/etc/s6-overlay/s6-rc.d/init-mod-healthchecks-apprise-add-package/run
#!/usr/bin/with-contenv bash

echo '
**************************************************************************
**************************************************************************
**************************************************************************
****                                                                  ****
****                                                                  ****
****          healthchecks-apprise mod has been deprecated            ****
****                                                                  ****
****                                                                  ****
****            apprise is now included in the main image             ****
****                                                                  ****
****                                                                  ****
**************************************************************************
**************************************************************************
**************************************************************************
'
