# root/etc/s6-overlay/s6-rc.d/init-mod-transmission-env-var-settings/up
/etc/s6-overlay/s6-rc.d/init-mod-transmission-env-var-settings/run
