# root/etc/s6-overlay/s6-rc.d/init-mod-universal-apprise-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-universal-apprise-add-package/run