# root/etc/s6-overlay/s6-rc.d/init-mod-swag-auto-proxy-setup/finish
#!/usr/bin/with-contenv bash

sed -i '/\/app\/auto-proxy.sh/d' /etc/crontabs/root
rm -rf /etc/nginx/http.d/auto-proxy*.conf /config/nginx/proxy-confs/auto-proxy*.conf
