# root/etc/s6-overlay/s6-rc.d/init-mod-plex-remove-codecs/run
#!/usr/bin/with-contenv bash 

echo "Remove codecs: start"

find "/config/Library/Application Support/Plex Media Server/Codecs/"* -type d -exec rm -rv {} +

echo "Remove codecs: end"
