# root/etc/s6-overlay/s6-rc.d/init-mod-plex-remove-codecs/up
/etc/s6-overlay/s6-rc.d/init-mod-plex-remove-codecs/run