# root/etc/s6-overlay/s6-rc.d/init-mod-universal-docker-in-docker-qemu/run
#!/usr/bin/with-contenv bash

echo "**** Enabling QEMU ****"
docker pull lscr.io/linuxserver/qemu-static
exec docker run --rm --privileged lscr.io/linuxserver/qemu-static --reset -p yes
