# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-prolog-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-code-server-prolog-add-package/run