# root/etc/s6-overlay/s6-rc.d/init-mod-daapd-snapclient/run
#!/usr/bin/with-contenv bash

# Alpine
if [ -f /sbin/apk ]; then
    echo "snapcast-client" >> /mod-repo-packages-to-install.list
fi
