# root/etc/s6-overlay/s6-rc.d/init-mod-daapd-snapclient/up
/etc/s6-overlay/s6-rc.d/init-mod-daapd-snapclient/run
