# root/etc/s6-overlay/s6-rc.d/init-mod-jellyfin-rffmpeg-setup/up
/etc/s6-overlay/s6-rc.d/init-mod-jellyfin-rffmpeg-setup/run