# root/etc/s6-overlay/s6-rc.d/init-mod-beets-drop2beets-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-beets-drop2beets-add-package/run
