# root/etc/s6-overlay/s6-rc.d/init-mod-lazylibrarian-ffmpeg-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-lazylibrarian-ffmpeg-add-package/run