# root/etc/s6-overlay/s6-rc.d/init-mod-plex-audnexus-add-package/run
#!/usr/bin/with-contenv bash 

if [ -z `command -v git` ]; then
    echo "**** adding plex-audnexus deps to package install list ****"
    echo "git" >> /mod-repo-packages-to-install.list
else
    echo "**** plex-audnexus deps already installed, skipping ****"
fi