# root/etc/s6-overlay/s6-rc.d/init-mod-nextcloud-memories-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-nextcloud-memories-add-package/run