# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-pnpm-install/run
#!/usr/bin/with-contenv bash

echo "**** installing pnpm  ****"
wget -qO- https://get.pnpm.io/install.sh | ENV="$HOME/.bashrc" SHELL="$(which bash)" bash -