# root/etc/s6-overlay/s6-rc.d/init-mod-emby-mediainfo-plugin-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-emby-mediainfo-plugin-add-package/run