# root/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc cloudflared tunnel --no-autoupdate --config /etc/cloudflared/config.yml run