# root/etc/s6-overlay/s6-rc.d/svc-mod-beets-httpshell/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc python3 /usr/local/bin/beets-httpshell.py
