# root/etc/s6-overlay/s6-rc.d/init-mod-swag-dashboard-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-swag-dashboard-add-package/run