# root/etc/s6-overlay/s6-rc.d/init-mod-swag-ffmpeg-add-package/run
#!/usr/bin/with-contenv bash

if ! apk info 2>&1 | grep -q "ffmpeg"; then
    echo "**** adding ffmpeg to package install list ****"
    echo "ffmpeg" >> /mod-repo-packages-to-install.list
else
    echo "**** ffmpeg already installed, skipping ****"
fi
