# root/etc/s6-overlay/s6-rc.d/svc-mod-daapd-snapclient/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    snapclient --player file:filename=/music/snapclient,mode=w
