# root/etc/s6-overlay/s6-rc.d/init-mod-radarr-striptracks-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-radarr-striptracks-add-package/run