# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-svn-install/run
#!/usr/bin/with-contenv bash

install-extension johnstoncode.svn-scm
