# root/etc/services.d/stdout-logs/run
#!/usr/bin/with-contenv bash

$TAIL_LOGS=$(echo $LOGS_TO_STDOUT | sed 's#|# -f #g')

echo "Executing: tail -f $TAIL_LOGS"
tail -f $TAIL_LOGS
