# root/etc/services.d/stdout-logs/run
#!/usr/bin/with-contenv bash

tail -f $LOGS_TO_STDOUT
