# root/etc/s6-overlay/s6-rc.d/init-mod-swag-maxmind-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-swag-maxmind-add-package/run