# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-terraform-install/run
#!/usr/bin/with-contenv bash

install-extension hashicorp.terraform
