# root/etc/s6-overlay/s6-rc.d/init-mod-firefox-fonts-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-firefox-fonts-add-package/run