# root/etc/s6-overlay/s6-rc.d/init-mod-netbox-slurpit-install/up
/etc/s6-overlay/s6-rc.d/init-mod-netbox-slurpit-install/run