# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-java11-add-package/run
#!/usr/bin/with-contenv bash

echo "**** Adding java11 packages to install list ****"
echo "openjdk-11-jdk" >> /mod-repo-packages-to-install.list
