# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-stdout-logs/finish
#!/usr/bin/with-contenv bash

if [ -z "${LOGS_TO_STDOUT}" ]; then
    exit 125
fi
