# root/etc/s6-overlay/s6-rc.d/init-mod-beets-httpshell/run
#!/usr/bin/with-contenv bash

echo "**** installing beets-httpshell mod ****"
echo "**** httpshell port: ${HTTPSHELL_PORT:-5555} ****"
echo "**** beets-httpshell mod installed ****"
