# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-nodejs-add-package/run
#!/usr/bin/with-contenv bash

echo "**** installing nodejs dev environment ****"
if ! dpkg -l | grep gnupg > /dev/null; then
    apt-get update && apt-get install -y --no-install-recommends gnupg
fi

echo "*** Downloading and installing nodejs package ***"
curl -fsSL https://deb.nodesource.com/setup_19.x | sudo -E bash - && sudo apt-get install -y nodejs