# root/etc/s6-overlay/s6-rc.d/init-mod-swag-dashboard-preload-logs/up
/etc/s6-overlay/s6-rc.d/init-mod-swag-dashboard-preload-logs/run