# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-flutter-add-package/run
#!/usr/bin/with-contenv bash

if ! which unzip >/dev/null; then
    echo "**** Adding flutter dependency packages to install list ****"
    echo "unzip" >> /mod-repo-packages-to-install.list
fi
