# root/etc/s6-overlay/s6-rc.d/init-mods-tshoot/up
/etc/s6-overlay/s6-rc.d/init-mods-tshoot/run
