# root/etc/s6-overlay/s6-rc.d/mod-init-tshoot/up
/etc/s6-overlay/s6-rc.d/mod-init-tshoot/run
