# root/etc/s6-overlay/s6-rc.d/init-mod-openssh-server-trusted-ca-install/up
/etc/s6-overlay/s6-rc.d/init-mod-openssh-server-trusted-ca-install/run
