# root/etc/s6-overlay/s6-rc.d/init-mod-openssh-server-autossh-add-package/run
#!/usr/bin/with-contenv bash

if [ -z `command -v autossh` ]; then
    echo "**** Adding autossh to package install list ****"
    echo "autossh" >> /mod-repo-packages-to-install.list
else
    echo "**** autossh already installed ****"
fi
