# root/etc/s6-overlay/s6-rc.d/init-mod-calibre-web-kcc-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-calibre-web-kcc-add-package/run