# root/etc/s6-overlay/s6-rc.d/init-mod-imagename-modname-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-imagename-modname-add-package/run