# root/etc/s6-overlay/s6-rc.d/init-mod-universal-infisical-cli-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-universal-infisical-cli-add-package/run