# root/etc/s6-overlay/s6-rc.d/init-mod-sabnzbd-par2cmdline-turbo-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-sabnzbd-par2cmdline-turbo-add-package/run
