# root/etc/services.d/geoip2influx/run
#!/usr/bin/with-contenv bash

  exec \
        python3 /config/geoip2db/geoip2influx.py