# root/etc/s6-overlay/s6-rc.d/svc-mod-beets-drop2beets/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ -z ${LSIO_NON_ROOT_USER} ]]; then
    exec \
        s6-setuidgid abc beet dropbox
else
    exec \
        beet dropbox
fi
