# root/etc/s6-overlay/s6-rc.d/init-mod-swag-imagemagick-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-swag-imagemagick-add-package/run