# root/etc/s6-overlay/s6-rc.d/svc-mod-swag-ondemand/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc python3 /app/swag-ondemand.py
